# .config/production/supervisor.conf
[program:uwsgi]
command = uwsgi --ini /srv/project/.config/production/uwsgi.ini

[program:frontend]
command = ng serve --host 0.0.0.0 --disable-host-check

[program:nginx]
command = nginx