# .config/production/supervisor.conf
[program:uwsgi]
command=uwsgi --ini /srv/project/.config/produciton/uwsgi.ini

[program:nginx]
command=nginx
